# gambler-trading.service
[Unit]
Description=GamblerAI Alpaca Paper Trading Bot
After=network.target

[Service]
Type=simple
User=ozay
WorkingDirectory=/home/ozay/GamblerAi
Environment="PATH=/home/ozay/GamblerAi/venv/bin"
ExecStart=/home/ozay/GamblerAi/venv/bin/python scripts/alpaca_paper_trading.py --continuous --symbols AAPL,MSFT,GOOGL,AMZN,TSLA,NVDA,META,AMD,NFLX,SPY,ADBE,CSCO,INTC,QCOM,AVGO,TXN,ASML,ORCL,PYPL,CMCSA,PEP,COST,SBUX,ABNB,GILD --interval 60
Restart=always
RestartSec=10
StandardOutput=append:/home/ozay/GamblerAi/logs/trading.log
StandardError=append:/home/ozay/GamblerAi/logs/trading.error.log

[Install]
WantedBy=multi-user.target
